# shutdownButton.service
[Unit]
Description=Shutdown Daemon

[Service]
ExecStart =/home/pi/shutdownButton.py
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
